# scripts/gpu-monitor.service
[Unit]
Description=GPU Monitor Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/gpu_monitor
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target 